# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:n8n-mcp-upstream]
command=/bin/sh -c 'set -a; . /app/upstream.env 2>/dev/null || exit 0; set +a; exec node /app/dist/mcp/index.js'
directory=/app
autostart=false
autorestart=true
startsecs=3
startretries=10
priority=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wrapper]
command=node /wrapper/dist/index.js
directory=/wrapper
autostart=true
autorestart=true
priority=2
environment=CONFIG_PATH="/app/wrapper-data/config.json",UPSTREAM_ENV_PATH="/app/upstream.env",PORT="9640",UPSTREAM_PORT="3000",DISCOVERY_PORT="9099"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
